# .claude/hooks/monitor-daemon/monitor-daemon.service
[Unit]
Description=Oppie CI/PR Monitor Daemon
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%WORKING_DIR%
Environment="GITHUB_TOKEN=%GITHUB_TOKEN%"
ExecStart=%INSTALL_DIR%/.claude/hooks/monitor-daemon/monitor-daemon
Restart=always
RestartSec=10
StandardOutput=append:/var/log/monitor-daemon.log
StandardError=append:/var/log/monitor-daemon.error.log

[Install]
WantedBy=multi-user.target